# desktop/configs/autostart/xfce.sh.desktop
[Desktop Entry]
Exec=/home/klein/xfce.sh
Icon=dialog-scripts
Name=xfce.sh
Path=
Type=Application
X-KDE-AutostartScript=true
